# crontab
* * * * * php /var/radius/interim.php
00 00 * * * sh /var/radius/clean.sh
